# fs-build/service/change_mac.service
[Unit]
Description=Expand partition size
After=sysinit.target

[Service]
Type=oneshot
ExecStart=/usr/lib/walnutpi/change_mac
RemainAfterExit=yes
StandardOutput=null

[Install]
WantedBy=multi-user.target
